# etc/udev/rules.d/73-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/bin/su -l arts /home/arts/bin/monitor-hotplug.sh"
